# TMR2026/systemd/carrito_tmr.service
[Unit]
Description=Carrito TMR 2026 — Sistema Autonomo
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/Carrito/TMR2026
Environment="PYTHONUNBUFFERED=1"
Environment="SDL_VIDEODRIVER=dummy"
Environment="SDL_AUDIODRIVER=dummy"
ExecStart=/usr/bin/python3 /home/pi/Carrito/TMR2026/main.py
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=carrito_tmr

[Install]
WantedBy=multi-user.target
